# scripts/medical-inventory.service
# TEMPLATE FILE - FOR REFERENCE ONLY
# This file is not used directly. The install_autostart.sh script
# dynamically generates the actual service file with correct paths.
# 
# To install: Run sudo ./install_autostart.sh
#
# This template shows the structure of the systemd service file
# that will be created at /etc/systemd/system/medical-inventory@.service

[Unit]
Description=Medical Inventory System for NASA Hunch
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
User=%i
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/%i/.Xauthority
WorkingDirectory=/home/%i/Medical-Inventory-System-NASA-Hunch-25-26
ExecStart=/home/%i/Medical-Inventory-System-NASA-Hunch-25-26/start_medical_inventory.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target
